# deploy/fortress-discord.service
[Unit]
Description=Fortress Command Discord Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=deploy
Group=deploy
WorkingDirectory=/home/deploy/OPENAGENT/discord_bot
EnvironmentFile=/home/deploy/.env.openagent
ExecStart=/home/deploy/OPENAGENT/discord_bot/.venv/bin/python3 /home/deploy/OPENAGENT/discord_bot/bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
